# python/uwsgi_supervisor.conf
[program:$APPLICATION_NAME]
command = $PATH_TO_UWSGI $PATH_TO_PYRAMID_CONFIG
user = $USER_NAME
environment=LANG=ru_RU.UTF-8
